# scripts/com.pathlock.metabase-daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pathlock.metabase-daily-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/havivrosh/work/metabase_kpi_sync/scripts/run_daily_sync.sh</string>
  </array>

  <!-- Run every day at 06:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/havivrosh/work/metabase_kpi_sync</string>

  <key>StandardOutPath</key>
  <string>/Users/havivrosh/work/metabase_kpi_sync/logs/launchd-daily-sync.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/havivrosh/work/metabase_kpi_sync/logs/launchd-daily-sync.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
